# deploy/systemd/invest-evolution.service
[Unit]
Description=Invest Evolution Web API
After=network.target
Wants=network.target

[Service]
Type=simple
User=invest
Group=invest
WorkingDirectory=/opt/invest-evolution
EnvironmentFile=-/etc/default/invest-evolution
ExecStart=/opt/invest-evolution/.venv/bin/gunicorn -c /opt/invest-evolution/gunicorn.conf.py wsgi:app
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
